# crontab
0 */2 * * * /path/to/backup.sh >> /var/log/mc_backup.log 2>&1