# crontab
0 */2 * * * cd /home/azureuser/mc && ./backup.sh >> /var/log/mc_backup.log 2>&1